# canopy.service
[Unit]
Description=canopy uwsgi microservices server

[Service]
WorkingDirectory=/var/orbit/cano.py/
ExecStart=/var/orbit/cano.py/start_server.sh
